# usb_defense/99-usb-monitor.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/local/bin/usb-monitor.sh"
